# deploy/discord-bot.service
[Unit]
Description=Discord-Discourse Bridge Bot
Documentation=https://github.com/ChrisTitusTech/community
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Set this to the non-root user that owns the plugin files on the *host*.
# The 'discourse' user only exists inside the Docker container — it does NOT
# exist on the host OS. Run: id -un   to find your current host username.
# Example: User=ubuntu  or  User=admin  or  User=deploy
# If you are unsure, run the bot as root by removing these two lines entirely.
User=root
# Group=  # uncomment and set if running as a non-root user

WorkingDirectory=/var/discourse/community/scripts/discord_bot

# Load secrets from the .env file — never commit this file to version control.
EnvironmentFile=/var/discourse/community/scripts/discord_bot/.env

ExecStart=/var/discourse/community/scripts/discord_bot/venv/bin/python3 bot.py

# Restart the bot if it crashes; wait 5 s between attempts.
Restart=on-failure
RestartSec=5

# Give Python 10 s to shut down gracefully before SIGKILL.
TimeoutStopSec=10

# Send stdout/stderr to the systemd journal (view with: journalctl -u discord-bot -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=discord-bridge-bot

[Install]
WantedBy=multi-user.target
